# devops/systemd/buddy-bringup(_IP).service
[Unit]
Description=Buddy ROS 2 bringup
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=robot
Environment=ROS_DOMAIN_ID=42
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
ExecStart=/home/robot/robot_ws/install/buddy_bringup/share/buddy_bringup/scripts/start_robot.sh
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
